# k8s/usr/lib/systemd/system/kubelet.service.d/10-config.conf
[Unit]
Requires=usr-libexec.mount
After=usr-libexec.mount

[Service]
EnvironmentFile=-/var/lib/kubelet/kubeadm-flags.env
EnvironmentFile=-/var/lib/kubelet/kube-flags.env
ExecStart=
ExecStart=/usr/bin/kubelet \
	    $KUBE_LOGTOSTDERR \
	    $KUBE_LOG_LEVEL \
	    $KUBELET_API_SERVER \
	    $KUBELET_ADDRESS \
	    $KUBELET_PORT \
	    $KUBELET_HOSTNAME \
	    $KUBE_ALLOW_PRIV \
	    $KUBELET_ARGS \
      --kubeconfig=/etc/kubernetes/kubelet.conf \
      --config=/var/lib/kubelet/config.yaml  \
			--bootstrap-kubeconfig=/etc/kubernetes/bootstrap-kubelet.conf \
      $KUBELET_KUBEADM_ARGS
