# mkosi.extra/build/k3s/files/usr/lib/systemd/system/k3s.service.d/10-agent.conf
[Unit]
After=docker.service
[Service]
ExecStart=
ExecStart=/usr/bin/k3s agent --docker --server https://172.16.1.20:6443 --token K1059b31b98ada5cf809a98486d100709333816aa2726ef06790af33fbe05145c3a::node:4587668b2aa3ac4c83e9bc4a0458e64c --data-dir /mnt/data/k3s
