# k8s/usr/lib/systemd/system/docker.service.d/10-storaget.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --data-root /mnt/data/docker
